# teletext-update.service
[Unit]
Description=Automatically check active teletext service for updates

[Service]
Type=oneshot
ExecStart=%h/.local/bin/teletext-update
