# vbit2.service
[Unit]
Description=VBIT2

[Service]
Type=simple
ExecStart=%h/vbit2/runvbit2.sh
Restart=always

[Install]
WantedBy=default.target
